# source_cfg/HQ-R/etc/systemd/system/etc_bckp.timer
[Unit]
Description=Backup /etc configs every day

[Timer]
OnBootSec=1min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target