# rootfs/BR-R/etc/systemd/system/etc_bckp.service
[Unit]
Description=backup /etc configs

[Service]
Type=oneshot
ExecStart=/root/etc_bckp.sh
